# client/heartbeat.service
[Unit]
Description=Heartbeat Client - 家里有电检测
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/heartbeat.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
